# linux_enhancer/systemd/linux_enhancer.timer
[Unit]
Description=Run Linux Enhancer Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
